# gunicorn.conf
[program:gunicorn]
command=gunicorn --workers=3 run:app -b localhost:5000
directory=/deploy
